# roundcube-catchall/rootfs/etc/s6-overlay/s6-rc.d/catchall-watchdog/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# Lightweight watchdog so HA sees the add-on as running. HA restarts the
# whole service tree on config changes, which re-runs init-catchall.

bashio::log.info "Catch-all plugin provisioned; idling."
exec sleep infinity
