# roundcube/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Roundcube Webmail
# Runs PHP-FPM
# ==============================================================================

bashio::log.info "Starting PHP-FPM..."
exec php-fpm84 -R --nodaemonize
