# roundcube-catchall/rootfs/etc/s6-overlay/s6-rc.d/init-catchall/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Roundcube Catch-all
# Provisions the catch-all Roundcube plugin into the shared plugin folder so
# the Roundcube Webmail add-on can pick it up.
# ==============================================================================

declare domain
declare identity_autocreate_raw
declare autologin_raw
declare autologin_user
declare autologin_password
declare fe_api_key
declare fe_auto_delete_raw
declare identity_autocreate_php
declare autologin_php
declare fe_auto_delete_php
declare domain_php
declare autologin_user_php
declare autologin_password_php
declare fe_api_key_php
declare -r plugin_repo="https://github.com/teh-hippo/roundcube-catchall.git"
# Pin to a specific commit to avoid silently pulling new upstream code into
# Roundcube every restart. Bump intentionally when rolling forward.
declare -r plugin_ref="7304fde131aa95088db785bb1c639fe96cc918f2"
declare -r plugin_root="/share/roundcube/plugins"
declare -r plugin_dir="${plugin_root}/catchall"
declare -r plugin_staging="${plugin_root}/.catchall.staging.$$"
declare -r plugin_previous="${plugin_root}/.catchall.previous"
declare -r marker_name=".ha-plugin-ready"

cleanup() {
    rm -rf "${plugin_staging}" 2>/dev/null || true
}
trap cleanup EXIT

bashio::log.info "Provisioning catch-all plugin..."

if ! bashio::config.has_value 'domain'; then
    bashio::log.fatal "domain is required."
    exit 1
fi

domain=$(bashio::config 'domain')
identity_autocreate_raw=$(bashio::config 'identity_autocreate')
autologin_raw=$(bashio::config 'autologin')
autologin_user=$(bashio::config 'autologin_user')
autologin_password=$(bashio::config 'autologin_password')
fe_api_key=$(bashio::config 'forward_email_api_key')
fe_auto_delete_raw=$(bashio::config 'forward_email_auto_delete')

if [[ "${domain}" == *"@"* ]]; then
    bashio::log.fatal "domain must be a bare domain (e.g. 'example.com'), not an email address."
    bashio::log.fatal "  got: '${domain}'"
    exit 1
fi
if [[ "${domain}" != *.* ]]; then
    bashio::log.fatal "domain must be a fully-qualified domain (contain a dot), got: '${domain}'"
    exit 1
fi

if [[ "${autologin_raw}" == "true" ]]; then
    if [[ -z "${autologin_user}" || -z "${autologin_password}" ]]; then
        bashio::log.fatal "autologin is enabled but autologin_user and/or autologin_password are empty."
        bashio::log.fatal "  Set both, or disable autologin."
        exit 1
    fi
    if [[ "${autologin_user}" != *"@${domain}" ]]; then
        bashio::log.warning "autologin_user '${autologin_user}' is not on domain '${domain}'."
    fi
    autologin_php="true"
else
    autologin_php="false"
fi

[[ "${identity_autocreate_raw}" == "false" ]] && identity_autocreate_php="false" || identity_autocreate_php="true"
[[ "${fe_auto_delete_raw}" == "true" ]] && fe_auto_delete_php="true" || fe_auto_delete_php="false"

# Escape a value so it can safely sit inside a PHP single-quoted string.
# PHP only interprets \\ and \' within single quotes.
php_escape_sq() {
    local s="$1"
    s="${s//\\/\\\\}"
    s="${s//\'/\\\'}"
    printf '%s' "$s"
}

mkdir -p "${plugin_root}"
rm -rf "${plugin_staging}"

# Fetch into a staging directory. Roundcube only loads plugins that carry a
# .ha-plugin-ready marker, so the partial tree is invisible until the swap.
bashio::log.info "Cloning plugin source (pinned ${plugin_ref:0:12})..."
if ! git clone --quiet "${plugin_repo}" "${plugin_staging}"; then
    bashio::log.fatal "Failed to clone ${plugin_repo}"
    exit 1
fi
if ! git -C "${plugin_staging}" checkout --quiet "${plugin_ref}"; then
    bashio::log.fatal "Failed to checkout ${plugin_ref}"
    exit 1
fi
rm -rf "${plugin_staging}/.git"

if [[ -f "${plugin_staging}/composer.json" ]]; then
    bashio::log.info "Installing plugin composer dependencies..."
    (cd "${plugin_staging}" && composer install --no-dev --no-interaction --optimize-autoloader 2>&1) \
        || bashio::log.warning "composer install reported errors."
fi

bashio::log.info "Writing plugin config..."
domain_php=$(php_escape_sq "${domain}")
autologin_user_php=$(php_escape_sq "${autologin_user}")
autologin_password_php=$(php_escape_sq "${autologin_password}")
fe_api_key_php=$(php_escape_sq "${fe_api_key}")

cat > "${plugin_staging}/config.inc.php" <<PHPEOF
<?php
// Managed by the Roundcube Catch-all add-on — edits will be overwritten.

\$config['catchall_domain']               = '${domain_php}';
\$config['catchall_identity_autocreate']  = ${identity_autocreate_php};
\$config['catchall_autologin']            = ${autologin_php};
\$config['catchall_autologin_user']       = '${autologin_user_php}';
\$config['catchall_autologin_pass']       = '${autologin_password_php}';

// Forward Email integration (optional). When the API key is empty the plugin
// stays in shared-credential mode and none of the FE-specific paths run.
\$config['catchall_fe_api_key_plain']     = '${fe_api_key_php}';
\$config['catchall_fe_auto_delete']       = ${fe_auto_delete_php};
PHPEOF
chmod 0640 "${plugin_staging}/config.inc.php"

# Write the readiness marker last so partial staging trees are never picked up.
touch "${plugin_staging}/${marker_name}"

# Atomically swap the live tree. rename(2) is atomic within a single filesystem,
# so the consumer side will only ever see the old tree or the fully-written new
# one, never a half-written mix.
rm -rf "${plugin_previous}"
if [[ -e "${plugin_dir}" ]]; then
    mv "${plugin_dir}" "${plugin_previous}"
fi
if ! mv "${plugin_staging}" "${plugin_dir}"; then
    bashio::log.fatal "Failed to activate plugin tree."
    # Attempt to restore previous tree so Roundcube is not left empty-handed.
    if [[ -e "${plugin_previous}" ]]; then
        mv "${plugin_previous}" "${plugin_dir}" || true
    fi
    exit 1
fi
rm -rf "${plugin_previous}"

# Clean up any leftover tree from the legacy forwardemail add-on so Roundcube
# does not accidentally load both plugins on the next restart.
rm -rf "${plugin_root}/forwardemail" "${plugin_root}/.forwardemail.previous" 2>/dev/null || true

fe_api_key_mask=""
if [[ -n "${fe_api_key}" ]]; then
    fe_api_key_mask="${fe_api_key:0:4}…${fe_api_key: -2}"
fi
bashio::log.info "Plugin provisioned at ${plugin_dir}"
bashio::log.info "  domain:              ${domain}"
bashio::log.info "  identity_autocreate: ${identity_autocreate_php}"
bashio::log.info "  autologin:           ${autologin_php}"
if [[ "${autologin_php}" == "true" ]]; then
    bashio::log.info "  autologin_user:      ${autologin_user}"
fi
if [[ -n "${fe_api_key_mask}" ]]; then
    bashio::log.info "  forward_email:       ${fe_api_key_mask} (per-alias API mode)"
    bashio::log.info "  fe_auto_delete:      ${fe_auto_delete_php}"
else
    bashio::log.info "  forward_email:       (none — shared credential mode)"
fi
bashio::log.notice "Restart the Roundcube Webmail add-on to load plugin changes."
