# roundcube/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Roundcube Webmail
# Configures nginx server blocks from templates
# ==============================================================================

bashio::log.info "Configuring nginx..."

# Generate Ingress server configuration
bashio::var.json \
    interface "$(bashio::addon.ip_address)" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf

# Generate direct access configuration if port 80 is mapped
if bashio::var.has_value "$(bashio::addon.port 80)"; then
    cp /etc/nginx/templates/direct.conf /etc/nginx/servers/direct.conf
fi

bashio::log.info "nginx configured."
