# roundcube-forwardemail/rootfs/etc/s6-overlay/s6-rc.d/forwardemail-watchdog/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# Lightweight watchdog so HA sees the add-on as running. HA restarts the
# whole service tree on config changes, which re-runs init-forwardemail.

bashio::log.info "Forward Email plugin provisioned; idling."
exec sleep infinity
